# etc/services.d/dnsmasq/run
#!/usr/bin/execlineb -P

/usr/sbin/dnsmasq --test --keep-in-foreground --user=dnsmasq --group=dnsmasq --conf-file=/etc/dnsmasq/dnsmasq.conf
